# deploy/sokolol.gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/home/o_adm/www/sokolol/deploy/sokolol.gunicorn.sock

[Install]
WantedBy=sockets.target